# install.service
[Unit]
Description=entrypoint install script

[Service]
ExecStart=/entrypoint.sh

[Install]
WantedBy=multi-user.target
